# etc/supervisor/conf.d/celery.conf
; ==================================
;  celery worker supervisor example
; ==================================

; the name of your supervisord program
[program:worker1]

; Set full path to celery program if using virtualenv
command=/home/vandm/Env/vandm/bin/celery -A vandm worker -l info
; The directory to your Django project
directory=/home/vandm/Env/vandm/vandm
; Put process stdout output in this file
stdout_logfile=/var/log/supervisor/celery/worker1/error.log

; Put process stderr output in this file
stderr_logfile=/var/log/supervisor/celery/worker1/out.log
startsecs=10

; If supervisord is run as the root user, switch users to this UNIX user account
; before doing any processing.
user=vandm

; Supervisor will start as many instances of this program as named by numprocs
numprocs=1

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600


; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; If true, this program will start automatically when supervisord is started
autostart=true

; May be one of false, unexpected, or true. If false, the process will never
; be autorestarted. If unexpected, the process will be restart when the program
; exits with an exit code that is not one of the exit codes associated with this
; process   ^`^y configuration (see exitcodes). If true, the process will be
; unconditionally restarted when it exits, without regard to its exit code.
autorestart=true


;[group:workers]
;programs=worker1,worker2,worker3
